# _gpt5intake/Users/sawyer/Library/LaunchAgents/com.thoughtmarks.watchdog.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thoughtmarks.watchdog.tunnel</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sawyer/gitSync/gpt-cursor-runner/scripts/watchdog-tunnel.sh</string>
        <string>start</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/sawyer/gitSync/gpt-cursor-runner</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/sawyer/gitSync/gpt-cursor-runner/logs/watchdogs/tunnel-watchdog-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/sawyer/gitSync/gpt-cursor-runner/logs/watchdogs/tunnel-watchdog-stderr.log</string>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>WATCHDOG_MODE</key>
        <string>launchd</string>
        <key>WATCHDOG_INSTALLATION_UUID</key>
        <string>4C4B566F-31D8-4B99-92BB-FC5BCC78661F</string>
    </dict>
</dict>
</plist>
