# _gpt5intake/Users/sawyer/Library/LaunchAgents/com.grammarly.ProjectLlama.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>GD_AGENT_LABEL</key>
		<string>com.grammarly.ProjectLlama.cleanup</string>
		<key>GD_AGENT_PLIST_PATH</key>
		<string>/Users/sawyer/Library/LaunchAgents/com.grammarly.ProjectLlama.cleanup.plist</string>
		<key>GD_AGENT_WATCH_PATHS</key>
		<string>/Users/sawyer/Applications:/Applications:/Network/Applications:/System/Cryptexes/App/System/Applications:/System/Cryptexes/OS/System/Applications:/System/Applications</string>
		<key>GD_BUNDLE_ID</key>
		<string>com.grammarly.ProjectLlama</string>
		<key>GD_BUNDLE_NAME</key>
		<string>Grammarly Desktop</string>
	</dict>
	<key>Label</key>
	<string>com.grammarly.ProjectLlama.cleanup</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/sawyer/Library/Application Support/com.grammarly.ProjectLlama/GRLoginHelper.app/Contents/MacOS/GRLoginHelper</string>
		<string>cleanup</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>WatchPaths</key>
	<array>
		<string>/Users/sawyer/Applications</string>
		<string>/Applications</string>
		<string>/Network/Applications</string>
		<string>/System/Cryptexes/App/System/Applications</string>
		<string>/System/Cryptexes/OS/System/Applications</string>
		<string>/System/Applications</string>
	</array>
</dict>
</plist>
