# _gpt5intake/Users/sawyer/Library/LaunchAgents/com.dropbox.DropboxUpdater.wake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AbandonProcessGroup</key>
	<true/>
	<key>AssociatedBundleIdentifiers</key>
	<array>
		<string>com.dropbox.DropboxUpdater</string>
	</array>
	<key>Label</key>
	<string>com.dropbox.DropboxUpdater.wake</string>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/sawyer/Library/Application Support/Dropbox/DropboxUpdater/Current/DropboxUpdater.app/Contents/MacOS/DropboxUpdater</string>
		<string>--wake-all</string>
		<string>--enable-logging</string>
		<string>--vmodule=*/components/update_client/*=2,*/chrome/updater/*=2</string>
	</array>
	<key>StartInterval</key>
	<integer>3600</integer>
</dict>
</plist>
