# _gpt5intake/Users/sawyer/Library/LaunchAgents/com.macpaw.CleanMyMac5.Updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AssociatedBundleIdentifiers</key>
	<string>com.macpaw.CleanMyMac5</string>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.macpaw.CleanMyMac5.Updater</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/sawyer/Library/Application Support/CleanMyMac_5/CleanMyMac_5_Updater.app/Contents/MacOS/CleanMyMac_5_Updater</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>21600</integer>
</dict>
</plist>
