# _gpt5intake/Users/sawyer/Library/LaunchAgents/com.thoughtmarks.cloudflared.gptcursorrunner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thoughtmarks.cloudflared.gptcursorrunner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--config</string>
    <string>/Users/sawyer/.cloudflared/gpt-cursor-runner-config.yml</string>
    <string>run</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/sawyer</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>NetworkState</key>
    <true/>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/sawyer/.cloudflared/gpt-cursor-runner.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sawyer/.cloudflared/gpt-cursor-runner.log</string>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StartInterval</key>
  <integer>60</integer>
</dict>
</plist>

