# _gpt5intake/Users/sawyer/Library/LaunchAgents/dev.sawyer.envsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.sawyer.envsync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/sawyer/gitSync/dev-tools/sync-env-daemon.js</string>
        <string>--watch</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/sawyer/gitSync/dev-tools/envsync.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/sawyer/gitSync/dev-tools/envsync.error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/sawyer/gitSync/dev-tools</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/sawyer</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ExitTimeOut</key>
    <integer>10</integer>
</dict>
</plist> 